# systemd/system-notifier.service
[Unit]
Description=Systemd Notifier - Telegram notifications for system events
Documentation=https://github.com/yourusername/systemd-notifier
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5

# Pass environment from systemd user session (critical for XDG_SESSION_ID)
PassEnvironment=XDG_SESSION_ID DBUS_SESSION_BUS_ADDRESS XDG_RUNTIME_DIR

# Environment variables
Environment="PYTHONUNBUFFERED=1"

# Working directory
WorkingDirectory=%h/.local/share/systemd-notifier

# Execute the notifier
ExecStart=%h/.local/share/systemd-notifier/src/notifier.py

# Standard output and error logging
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=128M
CPUQuota=10%

[Install]
WantedBy=default.target
